# shell/etc/supervisor/conf.d/ipmi-mqtt.conf
[program:ipmi-poll]
command=/usr/local/bin/ipmi-poll.sh

[program:mqtt-pub]
command=python3 /usr/local/bin/mqtt-pub.py
